# rootfs/etc/s6-overlay/s6-rc.d/symi-proxy/finish
#!/command/with-contenv bashio
# ==============================================================================
# Take down the S6 supervision tree when Symi Proxy fails
# ==============================================================================

declare exit_code
readonly exit_code_container=$(</run/s6-linux-init-container-results/exitcode)
readonly exit_code_service="$1"
readonly exit_code_signal="$2"

bashio::log.info \
  "Service exit with code ${exit_code_service}" \
  "(by signal ${exit_code_signal})"

if [[ "${exit_code_service}" -eq 256 ]]; then
  bashio::log.warning "Symi Proxy crashed, restarting..."
  exit 0
fi

bashio::log.warning "Symi Proxy failed, stopping add-on"
exec /run/s6/basedir/bin/halt 