# rootfs/etc/s6-overlay/s6-rc.d/symi-proxy/run
#!/command/with-contenv bashio
# ==============================================================================
# Start Symi Proxy service
# ==============================================================================

bashio::log.info "Starting Symi Proxy service..."
exec /run.sh 